# crontabs/portmgr
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.GlSF7IFIHx installed on Tue Nov 20 19:57:41 2012)
# (Cron version -- $FreeBSD: release/9.0.0/usr.sbin/cron/crontab/crontab.c 185041 2008-11-18 00:39:50Z matteo $)
#
# portmgr crontab
#
SCRIPTS="/var/portbuild/scripts/"
ERRORLOGS="/var/portbuild/errorlogs/"
#
0,15,30,45	*	*	*	*	cd ${ERRORLOGS} && ${SCRIPTS}/dopackagestats
