# etc/periodic/security/460.chkportsum.rej
***************
*** 42,61 ****
  
  case "${daily_status_security_chkportsum_enable}" in
  	[Yy][Ee][Ss])
  	pkg_info -ga 2>/dev/null | \
- 	while read one two three; do
- 		case ${one} in
  			Information)
- 			case ${two} in
- 				  for) name=${three%%:} ;;
- 					*) name='??' ;;
  			esac
  			;;
  			Mismatched|'') ;;
- 			*)
- 			if [ -n ${name} ]; then
- 				echo ${name}: ${one}
- 			fi
  			;;
  		esac
  	done
--- 15,34 ----
  
  case "${daily_status_security_chkportsum_enable}" in
  	[Yy][Ee][Ss])
+ 	set -f
  	pkg_info -ga 2>/dev/null | \
+ 	while IFS= read -r line; do
+ 		set -- $line 
+ 		case $1 in
  			Information)
+ 			case $2 in
+ 				for) name="${3%%:}" ;;
+ 				*) name='??' ;;
  			esac
  			;;
  			Mismatched|'') ;;
+ 			*) [ -n "${name}" ] && 
+ 				echo "${name}: ${line%% fails the original MD5 checksum}"
  			;;
  		esac
  	done
